# 99-garmin-fenix2.rules
KERNEL=="sdb", SUBSYSTEM=="block", ACTION=="add", ATTRS{idProduct}=="27af", ATTRS{idVendor}=="091e", RUN+="/bin/mount /dev/%k /mnt"
